# devops/supervisor/conf.d/010-nginx-service.ini
[program:nginx-service]
process_name=%(program_name)s_%(process_num)02d
priority=10
command=nginx -g "daemon off;"
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
stopwaitsecs=3600
killasgroup=true
stopasgroup=true