# devops/supervisor/conf.d/030-artisan-queue-service.ini
[program:artisan-queue-service-core]
process_name=%(program_name)s_%(process_num)02d
priority=30
command=php /var/www/artisan queue:work
autostart=true
autorestart=true
user=www-data
numprocs=1
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
stopwaitsecs=3600
killasgroup=true
stopasgroup=true